# taptap/rootfs/etc/s6-overlay/s6-rc.d/init-taptap/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Li Tin O've Weedle Add-on: taptap
# Configures TapTap to MQTT bridge
# ==============================================================================
declare taptap_log_level

if ! bashio::fs.directory_exists /run/taptap; then
    mkdir -p /run/taptap
fi

# Find the matching taptap log level
taptap_log_level="error"
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace|debug)
            taptap_log_level="debug"
            ;;
        info|notice)
            taptap_log_level="info"
            ;;
        warning)
            taptap_log_level="warning"
            ;;
        error|fatal|off)
            taptap_log_level="error"
            ;;
    esac
fi

if $(bashio::config.has_value 'taptap_modules'); then
    serials=""
    for serial in $(bashio::config 'taptap_modules'); do
        serials+=$serial","
    done
fi

# Generate taptap-mqtt configuration
bashio::var.json \
    mqtt_server "$(bashio::config 'mqtt_server')" \
    mqtt_port "$(bashio::config 'mqtt_port')" \
    mqtt_qos "$(bashio::config 'mqtt_qos')" \
    mqtt_timeout "$(bashio::config 'mqtt_timeout')" \
    mqtt_user "$(bashio::config 'mqtt_user')" \
    mqtt_pass "$(bashio::config 'mqtt_pass')" \
    taptap_log_level "${taptap_log_level}" \
    taptap_serial "$(bashio::config 'taptap_serial')" \
    taptap_address "$(bashio::config 'taptap_address')" \
    taptap_port "$(bashio::config 'taptap_port')" \
    taptap_modules "${serials%,}" \
    taptap_topic_prefix "$(bashio::config 'taptap_topic_prefix')" \
    taptap_topic_name "$(bashio::config 'taptap_topic_name')" \
    taptap_timeout "$(bashio::config 'taptap_timeout')" \
    taptap_update "$(bashio::config 'taptap_update')" \
    ha_discovery_prefix "$(bashio::config 'ha_discovery_prefix')" \
    ha_birth_topic "$(bashio::config 'ha_birth_topic')" \
    ha_nodes_availability_online "$(bashio::config 'ha_nodes_availability_online')" \
    ha_nodes_availability_identified "$(bashio::config 'ha_nodes_availability_identified')" \
    ha_strings_availability_online "$(bashio::config 'ha_strings_availability_online')" \
    ha_strings_availability_identified "$(bashio::config 'ha_strings_availability_identified')" \
    ha_stats_availability_online "$(bashio::config 'ha_stats_availability_online')" \
    ha_stats_availability_identified "$(bashio::config 'ha_stats_availability_identified')" \
    ha_nodes_sensors_recorder "$(bashio::config 'ha_nodes_sensors_recorder')" \
    ha_strings_sensors_recorder "$(bashio::config 'ha_strings_sensors_recorder')" \
    ha_stats_sensors_recorder "$(bashio::config 'ha_stats_sensors_recorder')" \
    | tempio \
        -template /etc/taptap/config.gtpl \
        -out /etc/taptap/config.ini

bashio::log.info "Generated taptap-mqtt configuration."