# ntopng/rootfs/etc/s6-overlay/s6-rc.d/geoip/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the netflow2ng daemon
# ==============================================================================

bashio::log.info "Starting Maxmind Database update runner..."

while true
do
    bashio::log.info "Geoip is configured"
    if bashio::fs.file_exists /etc/GeoIP.conf; then
        bashio::log.info "Checking state of the geoip database."
        cnt=$(find /data/geoip/*.mmdb -mtime -1 -print 2>/dev/null | wc -l)

        # Update GeoIP database if doesn't exist or is older then 1 day
        if [[ "$cnt" -eq 3 ]]; then
            bashio::log.info "Maxmind Database is up-to-date."
        else
            exec s6-setuidgid geoip geoip "rm -f /var/run/geoip/finished"
            bashio::log.info "Updating Maxmind Database..."
            exec s6-setuidgid geoip geoip /usr/bin/geoipupdate
            bashio::log.info "Maxmind Database updated."
            exec s6-setuidgid geoip geoip "touch /var/run/geoip/finished"
        fi
    fi
    sleep 21600
done
