# ntopng/rootfs/etc/s6-overlay/s6-rc.d/geoip/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the netflow2ng daemon
# ==============================================================================

bashio::log.info "Starting GeoIP database update runner..."

while true
do
    if bashio::fs.file_exists /etc/GeoIP.conf; then
        bashio::log.info "Checking state of the GeoIP database."
        if ! bashio::fs.directory_exists /data/geoip; then
            bashio::log.info "GeoIP database directory does not exist."
            sleep 5
            continue
        fi
        # Update GeoIP database if doesn't exist or is older then 1 day
        cnt=$(/usr/bin/find /data/geoip/*.mmdb -mtime -1 -print 2>/dev/null | /usr/bin/wc -l)
        bashio::log.info "GeoIP database $cnt files up-to-date."
        if [[ "$cnt" -eq 3 ]]; then
            bashio::log.info "GeoIP database is up-to-date."
        else
            rm -f /var/run/geoip/finished
            bashio::log.info "Updating GeoIP database..."
            exec /usr/bin/geoipupdate
            if [[ "$?" -ne 0 ]]; then
                bashio::log.error "GeoIP database update failed."
                sleep 60
                continue
            else
                bashio::log.info "GeoIP database updated."
                touch /var/run/geoip/finished
            fi
        fi
    else
        bashio::log.info "GeoIP configuration not found, skipping database update."
    fi
    sleep 21600
done
