# ntopng/rootfs/etc/s6-overlay/s6-rc.d/geoip/run
#!/command/with-contenv bashio
# ==============================================================================
# Li Tin O've Weedle Add-on: ntopng
# Runs the netflow2ng daemon
# ==============================================================================

bashio::log.info "Starting Maxmind Database update runner..."

while `true`
do
    if bashio::fs.file_exists /etc/GeoIP.conf; then
        cnt=$(find /data/geoip/*.mmdb -mtime -1 -print 2>/dev/null | wc -l)

        # Update GeoIP database if id doesn't exist or is older than 1 day
        if [[ "$cnt" -eq 3 ]]; then
            bashio::log.info "Maxmind Database is up-to-date."
        else
            rm -f /data/geoip/.finished
            bashio::log.info "Updating Maxmind Database..."
            /usr/bin/geoipupdate
            bashio::log.info "Maxmind Database updated."
            touch /data/geoip/.finished
        fi
    fi
    sleep 21600
done
